# etc/cron.daily/letsencrypt
#!/bin/sh
# Script that should renew certificates automatically.

[ -x /usr/bin/letsencrypt ] && /usr/bin/letsencrypt renew > /dev/nul
