# debian/fsinf-server-config.postinst
#!/bin/sh

set -e

# this should reload all sysctl rules, according to /etc/sysctl.d/README:
invoke-rc.d procps start

#DEBHELPER#
